# linux/cmds/kite-installer/rootless-system-data/.config/systemd/user/kite-updater.service.template
[Unit]
Description=Kite Updater

[Service]
Type=oneshot
StandardOutput=null
ExecStart={{.HomeDir}}/.local/share/kite/current/kite-update "update" "silent"
