# linux/cmds/kite-installer/rootless-system-data/.config/systemd/user/kite-updater.timer
[Unit]
Description=kite-updater timer

[Timer]
OnBootSec=5m
OnUnitInactiveSec=15m

[Install]
WantedBy=basic.target